# hr-bot-test/sapaedu.service

[Unit]
Description=SapaEdu Web Application
After=network.target postgresql.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/sapaedu
Environment=BOT_TOKEN=8317598561:AAG4kb5Dh1NxVTGLS60FIFy6_8UIHwh_RSo
Environment=GROUP_ID=-1002723413852
Environment=DATABASE_URL=postgresql://admin:Aa123456@10.128.0.79:5432/sapaedu
Environment=SECRET_KEY=your-secret-key-here-change-in-production
Environment=USE_FALLBACK_DATA=false
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 main.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
